# teclado_spanish.sh
#!/bin/bash
### BEGIN INIT INFO
# Provides: firewall
# Required-Start: $syslog $remote_fs $network
# Required-Stop: $syslog $remote_fs $network
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Description: Start firewall configuration
### END INIT INFO
setxkbmap -layout 'es,es' -model pc105
exit 0
